# setup/master.cron
* * * * * bash multi-nfsiostat/collect.bash df
1 0 * * * bash multi-nfsiostat/pull.bash yesterday clean
