# setup/slave.cron
* * * * * bash multi-nfsiostat/collect.bash
